# recipes-crypto/bitcoin/files/bitcoind.service
[Unit]
Description=Bitcoind service

[Service]
User=rockpro
ExecStart=bitcoind -disablewallet -printtoconsole -dbcache=100
